# work/ctwpro-rootless-5.6.0/final-audit-5.6.0-1/rootfs/var/jb/Library/LaunchDaemons/com.xxdevice.ctwsrv.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>RunAtLoad</key>
	<true/>
	<key>Label</key>
	<string>com.xxdevice.ctwsrv</string>
	<key>MachServices</key>
	<dict>
		<key>com.xxdevice.ctwsrv</key>
		<true/>
	</dict>
	<key>Program</key>
	<string>/var/jb/usr/bin/ctwsrv</string>
	<key>ProgramArguments</key>
	<array/>
	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
